# gpucheck/gpu-watchdog.timer
[Unit]
Description=Run GPU watchdog periodically

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min
Persistent=true

[Install]
WantedBy=timers.target
